# deploy/uswest/vpn-control.service
[Unit]
Description=Cudy VPN Control Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/cudy-control
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONUTF8=1
Environment=CUDY_ANDROID_APK_URL=https://github.com/Alesasha/Cudy_TR3000/releases/download/android-v1.46/NashVPN-CudyAgent-android-arm64-v1.46-20260721.apk
ExecStart=/usr/bin/python3 /opt/cudy-control/tools/vpn_control_app.py --db /opt/cudy-control/data/vpn_control.db --inventory /opt/cudy-control/config/vpn_inventory.json serve --host 127.0.0.1 --port 8765
Restart=always
RestartSec=5
User=cudy-control
Group=cudy-control
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/cudy-control

[Install]
WantedBy=multi-user.target
