# deploy/uswest/vpn-control-provider-refresh.timer
[Unit]
Description=Refresh Cudy VPN provider transports periodically

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true
Unit=vpn-control-provider-refresh.service

[Install]
WantedBy=timers.target
